# LakeMindControlCenter/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-out.log
stderr_logfile=/var/log/nginx-err.log

[program:bff]
command=uvicorn bff.app:app --host 127.0.0.1 --port 3001
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/bff-out.log
stderr_logfile=/var/log/bff-err.log

[program:steward]
command=uvicorn steward.app:app --host 127.0.0.1 --port 3002
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/steward-out.log
stderr_logfile=/var/log/steward-err.log
